# supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/sbin/mysqld

[program:jscity]
command=node /opt/jscity/js/server.js
